# deploy/sms_bark/sms-bark-forwarder.service
[Unit]
Description=Forward incoming SMS messages to Bark
After=network-online.target ModemManager.service
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /usr/local/bin/sms_forwarder.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
